# systemd/pureprivacy.service
[Unit]
Description=PurePrivacy Matrix appliance
Documentation=https://github.com/PurePrivacy/pureprivacy
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Edit WorkingDirectory to wherever you cloned the repo on this host.
WorkingDirectory=/opt/pureprivacy

# `After=docker.service` only sequences ordering — the docker socket may
# still be unresponsive when our unit fires (Docker forks a long-lived
# daemon and signals readiness before its API is fully serving).  Block
# until `docker info` actually succeeds, capped so a wedged docker
# doesn't hang the boot indefinitely.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do docker info >/dev/null 2>&1 && exit 0; sleep 2; done; echo "docker daemon not ready after 120s" >&2; exit 1'

ExecStart=/opt/pureprivacy/scripts/pureprivacy up
ExecStop=/opt/pureprivacy/scripts/pureprivacy stop

# Generous start timeout — first boot does image builds + initial Synapse
# DB migrations.  Subsequent boots finish in well under 60 s.
TimeoutStartSec=600

# Restart unconditionally if the up command itself exits non-zero.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
